# docker/systemd/coc-agent.service
[Unit]
Description=COC Runtime Agent
After=network.target
Documentation=https://github.com/chainofclaw/coc

[Service]
Type=simple
User=coc
Group=coc
WorkingDirectory=/opt/coc
ExecStart=/usr/bin/node --experimental-strip-types runtime/coc-agent.ts
Restart=always
RestartSec=10
LimitNOFILE=65535

Environment=COC_DATA_DIR=/var/lib/coc/runtime
Environment=COC_CONFIG=/etc/coc/runtime-agent.json
EnvironmentFile=-/etc/coc/coc-agent.env

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/coc /var/log/coc
PrivateTmp=true

StandardOutput=append:/var/log/coc/agent.log
StandardError=append:/var/log/coc/agent-error.log

[Install]
WantedBy=multi-user.target
